# release-monitor.service
[Unit]
Description=EreBus Release Monitor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/venv/bin/python __PROJECT_DIR__/monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
